# conf/supervisord.conf
[program:ihatemoney]
command=/path/to/your/app/venv/bin/gunicorn -c /etc/ihatemoney/gunicorn.conf.py ihatemoney.wsgi:application
user=www
autostart=true
autorestart=true
redirect_stderr=True
