# conf/supervisord.conf
[program:budget]
command=/path/to/your/app/venv/bin/gunicorn -c /path/to/your/app/conf/gunicorn.conf.py budget:app
directory=/path/to/your/app/budget/
user=www
autostart=true
autorestart=true
redirect_stderr=True
